# squirrel-capture.service
[Unit]
Description=Squirrel Soaker Camera Capture Daemon
After=network.target squirrel-trigger.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/squirrel_soaker
ExecStart=/usr/bin/python3 -u /home/pi/squirrel_soaker/capture.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
